# board/renmingrui/rootfs_overlay/usr/lib/systemd/system/Eplayer.service
[Unit]
Description=Eplayer startup
After=frp.service

[Service]
Type=simple
ExecStart=/root/Eplayer.sh

[Install]
WantedBy=multi-user.target
